# builder_dt/service/doozer_zkworker/run
#!/bin/sh
exec /bin/setuidgid chiseldata doozer run --interval 1 zkworker --loop
