# checkout_dt/service/doozer_checkout/run
#!/bin/sh
exec /bin/setuidgid chisel doozer run checkout
