# builder_dt/service/doozer_build/run
#!/bin/sh
exec /bin/setuidgid chiseldata doozer run overmind
