# reporter_puller/service/reporter_puller/run
#!/bin/sh
exec /bin/reporter_puller 2>&1
