# builder_dt/service/doozer_zkclean/run
#!/bin/sh
exec /bin/setuidgid chisel doozer run --interval 600 zkclean
